# build/tdb.conf
[program:tdb]
user = nobody
command = /bin/bash -c 'cd /usr/lib/tdb && /usr/bin/python2 ./lib/pex/TriggerDockerBuild.pex ./TriggerDockerBuild.py --config /config/config --logs /config/logs'
umask = 000
