# build/tdb.conf
[program:tdb]
user = nobody
command = cd /usr/lib/tdb && /usr/bin/python2 ./lib/pex/TriggerDockerBuild.pex ./TriggerDockerBuild.py
umask = 000
